# systemd-timer/atrack.service
[Unit]
Description=Runs ATrack
RefuseManualStart=no
RefuseManualStop=yes
 
[Service]
Type=oneshot
ExecStart=/home/ave/Documents/atrack.sh 